# src/templates/ethtool.conf.j2
[program:ethtool]
command=/sbin/ethtool -K $iface tx off sg off gro off gso off lro off tso off
auto_restart=false
startsecs=0
exitcodes=0,1
