# src/docker/suricata/supervisor.d/ethtool.conf
[program:ethtool]
command=/sbin/ethtool -K eth0 tx off sg off gro off gso off lro off tso off
auto_restart=false
startsecs=0
