# src/docker/suricata/supervisor.d/suricata.conf
[program:suricata]
command=/usr/bin/suricata
    -c /etc/suricata/suricata.yaml
    --pidfile /var/run/suricata.pid
    %(ENV_SURICATA_ARGS)s
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startsecs=10
